# .config/supervisor/deploy/supervisor-uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/app/bin/uwsgi --ini /srv/app/.config/uwsgi/deploy/uwsgi-app.ini