# .config/supervisor/debug/supervisor-uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/app/bin/uwsgi --ini /srv/app/.config/uwsgi/debug/uwsgi-app.ini